# deploy/supervisor-bit.conf
; Supervisor unit для БИТ.Технолог
; Положить в /etc/supervisor/conf.d/bit-technolog.conf
; Затем: supervisorctl reread && supervisorctl update && supervisorctl start bit-technolog

[program:bit-technolog]
command=/opt/beget/bit-technolog/serve.sh
directory=/opt/beget/bit-technolog
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stdout_logfile=/var/log/bit-technolog.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/bit-technolog.err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PYTHONUNBUFFERED="1",LANG="ru_RU.UTF-8",LC_ALL="ru_RU.UTF-8"
